# deploy/relay.service
[Unit]
Description=Crossworlds SSE Relay
After=network-online.target worldsim.service
Wants=network-online.target worldsim.service

[Service]
Type=simple
User=worldsim
Group=worldsim
ExecStart=/opt/worldsim/relay
Restart=always
RestartSec=5

EnvironmentFile=/opt/worldsim/relay.env

# Logging goes to journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay

# Security hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/worldsim

[Install]
WantedBy=multi-user.target
