# deploy/relay.service
[Unit]
Description=Crossworlds SSE Relay
After=network-online.target worldsim.service
Wants=network-online.target worldsim.service

[Service]
Type=simple
User=worldsim
Group=worldsim
ExecStart=/opt/worldsim/relay
Restart=always
RestartSec=5

EnvironmentFile=/opt/worldsim/relay.env

# Logging goes to journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay

# Security hardening (R83 brought to full parity with worldsim-backup-s3).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
ReadWritePaths=/opt/worldsim

[Install]
WantedBy=multi-user.target
