# deploy/sentinel.service
[Unit]
Description=Crossworlds Sentinel — Structural Health Monitor
After=worldsim.service
Requires=worldsim.service

[Service]
Type=simple
User=worldsim
Group=worldsim
WorkingDirectory=/opt/worldsim
ExecStart=/opt/worldsim/sentinel
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

# Security hardening (matches worldsim.service template).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/worldsim

[Install]
WantedBy=multi-user.target
