# deploy/worldsim-backup.timer
[Unit]
Description=Daily Crossworlds Database Backup

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
